# vibevm-gateway/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
user=root
priority=10

[program:vibevm-auth]
command=/usr/local/bin/uvicorn app:app --host 127.0.0.1 --port 8765
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/auth-stdout.log
stderr_logfile=/var/log/supervisor/auth-stderr.log
user=vibevm
priority=20
environment=HOME="/home/vibevm",USER="vibevm"
